# class/exercises/nginx/app.conf
description "daemon for node app"

start on startup
stop on shutdown
respawn

env PORT={{port}}

script
  {{node_path.stdout}} {{app_path}}
end script
