# class/exercises/services/app.conf
description "daemon for node app"

start on startup
stop on shutdown
respawn

script
  /usr/bin/node /mnt/app
end script
