# workshop/exercises/templates/solution/app.conf
description "daemon for node app"

start on startup
stop on shutdown
respawn

script
  {{node_path | default('/usr/bin/nodejs')}} {{app_path}}
end script
